# 72-security.rules
SUBSYSTEM=="usb", ENV{ID_VENDOR_ID}=="0e90", RUN+="/usr/local/bin/udev-lock"
